# shellinabox.conf
[supervisord]
nodaemon=true

[program:shellinabox]
command=__COMMAND__

